# mara-xporter.service
[Unit]
Description=mara-xporter

[Service]
ExecStart=/usr/local/bin/mara-xporter

[Install]
WantedBy=network.target

